# ops/com.greenwashing.monthly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
  <key>Label</key><string>com.greenwashing.monthly</string>
  <key>ProgramArguments</key><array>
    <string>__PROJECT_ROOT__/scripts/run_maintenance.sh</string><string>monthly</string>
  </array>
  <key>StartCalendarInterval</key><dict><key>Day</key><integer>1</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>30</integer></dict>
  <key>StandardOutPath</key><string>__PROJECT_ROOT__/.gw/logs/monthly.log</string>
  <key>StandardErrorPath</key><string>__PROJECT_ROOT__/.gw/logs/monthly-error.log</string>
</dict></plist>
